# supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:cron]
command=cron -f

[program:webserver]
command=python -m http.server 42069 --directory /app
